# systemd/pg-backup.timer
[Unit]
Description=Daily PostgreSQL backup to S3

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
