# systemd/pg-backup.timer
[Unit]
Description=Weekly PostgreSQL backup (Sunday 02:00, day before scraper)

[Timer]
OnCalendar=Sun *-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
