# services/disk-alert/disk-alert.service
[Unit]
Description=Disk usage alert via Telegram
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/push-monitor/disk-alert.env
ExecStart=/home/victor/infra/scripts/disk-alert.sh
User=victor
